# examples/docker/varnish/files/etc/supervisor/conf_enabled/log.conf
[program:log]
command=varnishncsa -g request -q 'ReqMethod eq "PURGE"'
autostart=true
autorestart=true
stopsignal=QUIT
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
